# debian/gitolite-backup-scripts.cron.hourly
#!/bin/sh

. /etc/gitolite-backup-scripts.conf

su $GITOLITEUSER -c "/usr/lib/gitolite-backup-scripts/bin/gitolite_backup_run" 
