# debian/gitolite-backup-scripts.cron.daily
#!/bin/sh

. /etc/gitolite-backup-scripts.conf

su $GITOLITEUSER -c "/usr/bin/gitolite_backup_run" 
